# deploy/horus.service.example
[Unit]
Description=Horus Django app
After=network.target

[Service]
User=deploy
Group=www-data
WorkingDirectory=/var/www/horus
EnvironmentFile=/var/www/horus/.env
ExecStart=/var/www/horus/venv/bin/gunicorn config.wsgi:application --bind 127.0.0.1:8003 --workers 3
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
